# osd-0.service
[Unit]
Description=OSD start Service
After=syslog.target

[Service]
Type=simple
ExecStart= 
Restart=always

[Install]
WantedBy=multi-user.target
